# PostgreSQL/Ansible/roles/wal-g/templates/walg-fullbackup.timer.j2
[Unit]
Description=Run walg-fullbackup.service every day
After=network.target
[Timer]
OnCalendar={{walg_full_backup_timer}}
[Install]
WantedBy=multi-user.target